# supervisior.conf
[include]
files=/etc/supervisor/conf.d/fun_supervisor.conf

[program:awesome_app]
directory=/home/monk/fun/fun
command=gunicorn -w 4 main:app