# supervisior.conf
[include]
files=/etc/supervisord.conf

[program:awesome_app]
directory=/home/monk/fun/fun
command=gunicorn -w 4 main:app